# lib/systemd/system/kdm.service.d/30_rads.conf
## Copyright (C) 2016 - 2022 ENCRYPTED SUPPORT LP <adrelanos@whonix.org>
## See the file COPYING for copying conditions.

[Unit]
Before=rads.service
After=rads-block-display-manger.service
ConditionPathExists=!/run/rads/rads-block-display-manger-startup.status
